# .config/systemd/user/automatic-git-fetch.service
[Unit]
Description=Update Git repositories

[Service]
Type=oneshot
ExecStart=/bin/sh /home/coates/bin/automatic-git-fetch.sh
