# .config/systemd/user/automatic-git-fetch.timer
[Unit]
Description=Timer to update Git repositories

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
